# stages/05-Wifibroadcast/FILES/overlay/etc/systemd/system/qopenhd.service
[Unit]
Description=QOpenHD
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/QOpenHD

[Install]
WantedBy=multi-user.target
